# deploy/systemd/styles4dogs-notification.timer
[Unit]
Description=Process Styling 4 Dogs notification queue and booking housekeeping every minute
Documentation=file:/opt/styles4dogs/share/NOTIFICATIONS.md

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true
Unit=styles4dogs-notification.service

[Install]
WantedBy=timers.target
