# deploy/systemd/styles4dogs-backup.timer
[Unit]
Description=Daily backup timer for Styling 4 Dogs
Documentation=file:/opt/styles4dogs/share/DEPLOYMENT.md

[Timer]
OnCalendar=*-*-* 03:20:00
RandomizedDelaySec=30m
Persistent=true
Unit=styles4dogs-backup.service

[Install]
WantedBy=timers.target
